# templates/startup_supervisor.conf.j2
[program:rsyslog_sentry]
command=python {{rsyslog_sentry_installdir}}/rsyslog_sentry.py
autostart=true
autorestart=true
